# sources/cattleshed.upstart.in
description "Backend of Wandbox"
start on runlevel [2345]:
stop on runlevel [!2345]:
respawn
chdir /
exec env - su cattleshed -c @@bindir@@/cattlegrid
