# sources/kennel.upstart.in
description "Frontend of Wandbox"
start on runlevel [2345]:
stop on runlevel [!2345]:
respawn
chdir /var/lib/kennel
exec env - su kennel -c '@@bindir@@/kennel Production'
